# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; nodiCore — REST API, internal only (not exposed outside the container)
[program:nodi-core]
command=dotnet /app/core/nodiCore.dll
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; CORE_PORT is set by start.sh from the container env var (default 5100)
environment=ASPNETCORE_URLS="http://+:%(ENV_CORE_PORT)s"

; nodiWeb — Blazor Server, this is the only port exposed to the host
[program:nodi-web]
command=dotnet /app/web/nodiWeb.dll
directory=/app/web
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Talks to nodi-core on localhost inside the same container
environment=ASPNETCORE_URLS="http://+:%(ENV_WEB_PORT)s",ApiBaseUrl="http://localhost:%(ENV_CORE_PORT)s"
